# docker/local/etc/profile.d/xdebug-php.sh
#!/usr/bin/bash

alias xdebug-php="XDEBUG_CONFIG=\"mode=debug client_host=192.168.90.1 client_port=9003 start_with_request=yes\" php"
alias xdebug-on="export XDEBUG_CONFIG=\"mode=debug client_host=192.168.90.1 client_port=9003 start_with_request=yes\""
alias xdebug-off="export XDEBUG_CONFIG="
